# scripts/launchd/com.shingihou.promote-corrections.plist
<?xml version="1.0" encoding="UTF-8"?>
<!DOCTYPE plist PUBLIC "-//Apple//DTD PLIST 1.0//EN"
  "http://www.apple.com/DTDs/PropertyList-1.0.dtd">
<plist version="1.0">
<dict>
  <key>Label</key>
  <string>com.shingihou.promote-corrections</string>
  <key>ProgramArguments</key>
  <array>
    <string>/Users/lin/voice-input/venv/bin/python3</string>
    <string>/Users/lin/voice-input/scripts/dictionary_promote_from_history.py</string>
    <string>--apply</string>
  </array>
  <key>StartCalendarInterval</key>
  <dict>
    <key>Hour</key>
    <integer>3</integer>
    <key>Minute</key>
    <integer>30</integer>
  </dict>
  <key>StandardOutPath</key>
  <string>/Users/lin/.voice-input/promote-corrections.stdout.log</string>
  <key>StandardErrorPath</key>
  <string>/Users/lin/.voice-input/promote-corrections.stderr.log</string>
  <key>EnvironmentVariables</key>
  <dict>
    <key>PATH</key>
    <string>/usr/local/bin:/usr/bin:/bin:/opt/homebrew/bin</string>
    <key>HOME</key>
    <string>/Users/lin</string>
  </dict>
</dict>
</plist>
